# root/app/epg/supervisord.conf
[program:epg]
command = /app/epg/xmltv.se.sh
autorestart = true
pidfile = /run/epg.pid
stdout_logfile = /var/log/epg.log
logfile_maxbytes = 50MB
priority = 2
