# root/app/epg/supervisord.conf
[program:epg]
command = /app/epg/update.sh --as-daemon
autorestart = true
pidfile = /run/epg.pid
stdout_logfile = /var/log/epg.log
logfile_maxbytes = 50MB
priority = 2
