# root/app/iptv/supervisord.conf
[program:iptv]
command = /app/iptv/update.sh --as-daemon
autorestart = true
pidfile = /run/iptv.pid
stdout_logfile = /var/log/iptv.log
logfile_maxbytes = 50MB
priority = 2
