# root/app/iptv/supervisord.template.conf
[program:{channel}]
command = /app/iptv/channel.sh {channel}
autorestart = true
pidfile = /run/{channel}.pid
stdout_logfile = /var/log/{channel}.log
logfile_maxbytes = 50MB
priority = 2

